# bin/league-pbe-pinger/pbeping.service
[Unit]
Description=Ping League PBE server and send discord messages when state changes
After=network.target

[Service]
Type=simple
User=hemul
WorkingDirectory=/home/hemul/dotfiles/bin/league-pbe-pinger
ExecStart=/home/hemul/dotfiles/bin/league-pbe-pinger/pbeping.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
