# bin/league-pbe-pinger/pbeping.service
[Unit]
Description=League PBE Status Pinger
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/hemul/dotfiles/bin/league-pbe-pinger
ExecStart=/home/hemul/dotfiles/bin/league-pbe-pinger/pbeping.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
